# supervisord.conf
# =============================================================================
# Supervisord Configuration for A.E.G.I.S. v8.0
# =============================================================================
# Manages three independent processes with auto-restart:
# - Feed Process: WebSocket data ingestion
# - Brain Process: ML analysis + trading signals
# - Orchestrator Process: Cache reconciliation + monitoring
# =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
loglevel=info

# Keep logs in stdout for Railway/Docker visibility
silent=false

# ============================================================================
# ORCHESTRATOR PROCESS - Monitoring + Maintenance (START FIRST!)
# ============================================================================
# Reconciliation, system monitoring, auto-maintenance
# PRIORITY 999 = Starts FIRST (before feed/brain)
# This ensures database + ring buffer initialization happens before feed/brain attach
# Restarts automatically if it crashes

[program:orchestrator]
command=python -m src.main orchestrator
directory=/home/runner/workspace
autostart=true
autorestart=true
startsecs=5
startretries=10
stopasgroup=true
killasgroup=true
priority=999

# Redirect logs to container stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Process environment
environment=PYTHONUNBUFFERED=1

# ============================================================================
# FEED PROCESS - Data Ingestion (START SECOND)
# ============================================================================
# Reads from Binance WebSocket and writes to ring buffer
# PRIORITY 100 = Starts after orchestrator
# Restarts automatically if it crashes

[program:feed]
command=python -m src.main feed
directory=/home/runner/workspace
autostart=true
autorestart=true
startsecs=5
startretries=10
stopasgroup=true
killasgroup=true
priority=100

# Redirect logs to container stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Process environment
environment=PYTHONUNBUFFERED=1

# ============================================================================
# BRAIN PROCESS - ML Analysis + Trading (START THIRD)
# ============================================================================
# Reads from ring buffer, runs analysis, executes trades
# PRIORITY 50 = Starts after feed
# Restarts automatically if it crashes

[program:brain]
command=python -m src.main brain
directory=/home/runner/workspace
autostart=true
autorestart=true
startsecs=5
startretries=10
stopasgroup=true
killasgroup=true
priority=50

# Redirect logs to container stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Process environment
environment=PYTHONUNBUFFERED=1

# ============================================================================
# Supervisord Process Group
# ============================================================================
# Manage all three processes as a group

[group:aegis]
programs=feed,brain,orchestrator
priority=999

# ============================================================================
# Configuration Notes:
# ============================================================================
#
# 1. PROCESS GROUPS:
#    - feed: Handles WebSocket connections and ring buffer writes
#    - brain: Processes data and executes trades
#    - orchestrator: System monitoring and maintenance
#
# 2. AUTO-RESTART:
#    - autostart=true: Start all processes when supervisord starts
#    - autorestart=true: Restart if any process exits abnormally
#    - startretries=10: Try restarting up to 10 times before giving up
#
# 3. LOGGING:
#    - stdout_logfile=/dev/stdout: Visible in container logs
#    - stderr_logfile=/dev/stderr: Visible in container logs
#    - No file buffering: maxbytes=0 for real-time visibility
#
# 4. GRACEFUL SHUTDOWN:
#    - stopasgroup=true: Send signals to entire process group
#    - killasgroup=true: Kill entire group if needed
#
# 5. RAILWAY/DOCKER:
#    - nodaemon=true: Keep supervisord in foreground
#    - PYTHONUNBUFFERED=1: Unbuffered Python output
#    - Logs redirect to stdout for Docker/Railway visibility
#
# =============================================================================
